# lib/process_messages.cron
#!/bin/bash

# /etc/crontab:
# 0 0 * * * cvzprj /home/cvzprj/mailbox_tester/current/lib/process_messages.cron

export LD_LIBRARY_PATH=/usr/local/lib
cd /home/cvzprj/mailbox_tester/current
bundle exec rake process_messages RAILS_ENV=production